# systemd/beatcheck-refresh.service
[Unit]
Description=BeatCheck Feed Refresh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/beatcheck --refresh
Restart=on-failure
RestartSec=5m
