# systemd/beatcheck-refresh.timer
[Unit]
Description=BeatCheck Feed Refresh Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
