# systemd/speedy-reader-refresh.service
[Unit]
Description=SpeedyReader Feed Refresh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/speedy-reader --refresh
Restart=on-failure
RestartSec=5m
